# deployment/amy_supervisor.conf
[program:amy-web]
user = amy
group = amy
command = /home/amy/amy/runserver.sh
autostart = true
stopasgroup=true
killasgroup=true
